# deploy/systemd/mopro-backup.service
[Unit]
Description=Mopro nightly Postgres backup to B2 + Hetzner via restic
Documentation=file:///opt/mopro/docs/runbooks/backup-failure.md
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=mopro
Group=mopro
# All secrets (RESTIC_PASSWORD, B2_*, HETZNER_*, HEALTHCHECK_*, SLACK_*) live here.
EnvironmentFile=/opt/mopro/.env

ExecStart=/opt/mopro/deploy/scripts/backup-postgres.sh

StandardOutput=journal
StandardError=journal
SyslogIdentifier=mopro-backup

# Resource guards — backup is a background job, not latency-sensitive.
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
CPUQuota=50%

# Retry once on failure (systemd triggers Restart; Healthchecks.io grace timer
# provides the final alert if both attempts fail).
Restart=on-failure
RestartSec=60
StartLimitIntervalSec=300
StartLimitBurst=2
