# deploy/systemd/mopro-backup.service
[Unit]
Description=Mopro daily Postgres backup to Hetzner Storage Box
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=mopro
Group=mopro
EnvironmentFile=/etc/mopro/.env
ExecStart=/opt/mopro/deploy/scripts/backup-postgres.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mopro-backup

# Resource guards
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
